# ticketing/msp-ticketing.service
[Unit]
Description=MSP Ticketing System
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/msp-ticketing
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
